# overlay_template/system/rootfs_overlay/usr/lib/systemd/system/8890bootpulse.service
[Unit]
Description=send boot pulses to cpu scheduler on user input

[Service]
Restart=always
ExecStart=/usr/sbin/8890boostpulse


