# overlay_template/system/rootfs_overlay/etc/init/8890boostpulse.conf
start on started lightdm

respawn

exec /usr/sbin/8890boostpulse
