# overlay_template/system/rootfs_overlay/etc/init/hide_sdcard.conf
start on startup

script
	if [ -n "$(df | grep -E '/dev/disk/by-partlabel/UT_.+_SD')" ]
	then
		mount -o bind /dev/null /lib/udev/rules.d/99-zzz-sdcard.rules
	fi
end script
