# overlay_template/system/rootfs_overlay/etc/init/limit_charging.conf
start on started lightdm

script
	if [ -e /home/phablet/.config/limit_charging ]
	then
		echo 1 > /sys/class/power_supply/battery/store_mode
	fi
end script
