# overlay_template/system/rootfs_overlay/etc/init/zram-8890.conf
start on started lightdm

script
	zramctl /dev/zram0 -s 2147483648
	mkswap /dev/zram0
	swapon /dev/zram0
end script
