# overlay_template/system/usr/lib/systemd/system/debug_telnet.service
[Unit]
Description=run debug telnet

[Service]
ExecStart=/usr/bin/bash /usr/libexec/lxc-android-config/device-hacks.d/debug_telnet.sh
Type=forking

